# units/user/cmdline.target.in

[Unit]
Description=Non-graphical services
Wants=xorg.target
Requires=dbus.socket
After=xorg.target dbus.socket

[Install]
Alias=default.target
