# units/user/e17.target

[Unit]
Description=The Englightenment Desktop Environment built on EFL
Wants=xorg.target
Requires=dbus.socket

