# units/user/enlightenment.service.in

[Unit]
Description=The enlightenment Window Manager
After=xorg.target
Requires=dbus.socket

[Service]
ExecStart=@EPREFIX@/bin/enlightenment_start

