# units/user/enlightenment.service.in

[Unit]
Description=The enlightenment Window Manager
After=xorg.target
Requires=dbus.socket

[Service]
ExecStart=@E_PROG@

[Install]
WantedBy=e17.target
