# units/user/enlightenment.service

[Unit]
Description=The enlightenment Window Manager
After=xorg.target
Requires=dbus.socket

[Service]
ExecStart=/usr/bin/enlightenment_start

