# units/user/enlightenment.service.in

[Unit]
Description=The enlightenment Window Manager
Requires=xorg.target
After=xorg.target

[Service]
ExecStart=@E_PROG@
Restart=on-failure

[Install]
Alias=window-manager.service
