# units/user/gnome.target.in

[Unit]
Description=The GNOME Desktop Environment
Wants=xorg.target
Requires=dbus.socket
