# units/user/gnome-session.service.in

[Unit]
Description=The GNOME session manager
After=xorg.target
Requires=dbus.socket

[Service]
ExecStart=@GNOME_PROG@
Restart=on-failure

[Install]
WantedBy=gnome.target
