# units/user/gnome-session.service.in

[Unit]
Description=The GNOME session manager
Requires=dislpay-server.target
After=display-server.target

[Service]
ExecStart=@GNOME_PROG@
Restart=on-failure

[Install]
Alias=session-manager.service
