# units/user/gnome-session.service.in

[Unit]
Description=The GNOME session manager
After=xorg.target
Requires=dbus.socket

[Service]
ExecStart=$GNOMEPREFIX@/bin/gnome-session

