# units/user/kde.service.in

[Unit]
Description=The KDE start agent
Requires=display-server.target
After=display-server.target

[Service]
ExecStart=@KDE_PROG@
Restart=on-failure

[Install]
Alias=session-manager.service
