# units/user/kde.target.in

[Unit]
Description=The KDE Desktop Environment
Wants=xorg.target
Requires=dbus.socket
AllowIsolate=true
