# units/user/mate-session.service.in

[Unit]
Description=The MATE session manager
Requires=xorg.target
After=xorg.target

[Service]
ExecStart=@MATE_PROG@
Restart=on-failure

[Install]
Alias=session-manager.service
