# units/user/mate.target.in

[Unit]
Description=The MATE Desktop Environment
Wants=xorg.target
Requires=dbus.socket
AllowIsolate=true
