# units/user/mythtv.target.in

[Unit]
Description=MythTV Frontend
Wants=xorg.target
Requires=dbus.socket
AllowIsolate=true

