# units/user/openbox.target.in

[Unit]
Description=A Desktop Environment built on Openbox
Wants=xorg.target
Requires=dbus.socket
AllowIsolate=true
