# units/user/openbox.service.in

[Unit]
Description=The Openbox Window Manager
Requires=xorg.target
After=xorg.target

[Service]
ExecStart=@OPENBOX_PROG@
Restart=on-failure

[Install]
Alias=window-manager.service
