# units/user/steam.target.in

[Unit]
Description=Steam fullscreen Desktop Environment
Wants=xorg.target
Requires=dbus.socket
AllowIsolate=true
