# units/user/steam.service.in

[Unit]
Description=Steam - Big Picture mode
Requires=display-server.target
After=display-server.target

[Service]
ExecStart=@STEAM_PROG@ steam://open/bigpicture
Restart=on-failure

[Install]
WantedBy=graphical.target
