# units/user/xbmc.target

[Unit]
Description=XBMC
Wants=xorg.target dbus.socket

