# units/user/xbmc.target.in

[Unit]
Description=XBMC
Wants=xorg.target
Requires=dbus.socket
AllowIsolate=true
