# units/user/xfce4-session.service.in

[Unit]
Description=Xfce4 Session Manager
Requires=display-server.target
After=display-server.target

[Service]
ExecStart=@XFCE_PROG@
Restart=on-failure

[Install]
Alias=session-manager.service

