# units/user/xfce4.target.in

[Unit]
Description=The Xfce 4.x Desktop Environment
Wants=xorg.target dbus.socket

