# units/user/xfce4-session.service.in

[Unit]
Description=Xfce4 Session Manager
After=xorg.target

[Service]
ExecStart=/usr/bin/xfce4-session

