# units/user/xfce4-session.service.in

[Unit]
Description=Xfce4 Session Manager
After=xorg.target

[Service]
ExecStart=@XFCE_PROG@
Restart=on-failure

[Install]
WantedBy=xfce4.target

