# units/user/xfwm4.service.in

[Unit]
Description=Xfce Window Manager Daemon
After=xorg.target

[Service]
ExecStart=@XFWM_PROG@
Restart=on-failure

[Install]
WantedBy=xfce4.target
