# units/user/xfwm4.service.in

[Unit]
Description=Xfce Window Manager Daemon 
After=xorg.target

[Service]
Restart=always
RestartSec=1
ExecStart=/usr/bin/xfwm4
