# units/user/xfwm4.service.in

[Unit]
Description=Xfce Window Manager Daemon
Requires=xorg.target
After=xorg.target

[Service]
ExecStart=@XFWM_PROG@
Restart=on-failure

[Install]
Alias=window-manager.service
