# units/system/systemd-stop-user-sessions.service.in
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Disable User Sessions
Documentation=man:systemd-user-sessions.service(8)
Before=shutdown.target
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@SYSTEMDUTILDIR@/systemd-user-sessions stop

[Install]
WantedBy=shutdown.target
