# units/system/user-session@.service.in
#  This file is part of systemd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU General Public License as published by
#  the Free Software Foundation; either version 2 of the License, or
#  (at your option) any later version.

[Unit]
Description=User Manager for %I
After=systemd-user-sessions.service
# This is needed for folks running Xorg on vt/1
Conflicts=getty@tty1.service

[Service]
User=%I
PAMName=login
# Note memory:/ required kernel memory CGroups support
ControlGroup=%R/user/%I/shared cpu:/ memory:/
ControlGroupModify=yes
Type=notify
TTYPath=/dev/tty1
ExecStart=-@SYSTEMDUTILDIR@/systemd --user
Environment=DISPLAY=:0
Environment=XDG_CACHE_HOME=%h/.cache
Environment=XDG_CONFIG_DIRS=/etc/xdg
Environment=XDG_CONFIG_HOME=%h/.config
Environment=XDG_DATA_DIRS=/usr/local/share/:/usr/share/
Environment=XDG_DATA_HOME=%h/.local/share
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/dbus/user_bus_socket
